# deps/redis.service
[Unit]
Description=redis

[Service]
ExecStart=/usr/bin/rkt run \
	--net=host \
        --volume redis,kind=host,source=/tmp/redisdata \
        unixvoid.com/redis

[Install]
WantedBy=multi-user.target
